# apriltag_detector.service
[Unit]
Description=Apriltag Detector Service

[Service]
Type=simple
ExecStart=cd %p && /usr/bin/python3 %s

[Install]
WantedBy=multi-user.target