# apriltag_detector.service
[Unit]
Description=Apriltag Detector Service

[Service]
Type=simple
WorkingDirectory = %P
ExecStart=/usr/bin/python3 %s

[Install]
WantedBy=multi-user.target