# apriltag_detector.service
[Unit]
Description=Apriltag Detector Service

[Service]
Type=simple
ExecStart=python3 %s

[Install]
WantedBy=multi-user.target