# apriltag_detector.service
[Unit]
Description=Dash Button Sniffer

[Service]
Type=simple
ExecStart=%s

[Install]
WantedBy=multi-user.target